# docker/all-in-one/rootfs/etc/s6-overlay/s6-rc.d/smtpd/run
#!/command/with-contenv sh
exec /smtpd serve -config /etc/infodancer/config.toml
