# docker/all-in-one/rootfs/etc/s6-overlay/s6-rc.d/webadmin/run
#!/command/with-contenv sh
exec /webadmin -config /etc/infodancer/config.toml
