# docker/all-in-one/rootfs/etc/s6-overlay/s6-rc.d/queue-manager/run
#!/command/with-contenv sh
exec /queue-manager --config /etc/infodancer/config.toml \
    --queue /var/spool/mail-queue --binary /mail-remote
